# nip05-dm-listener/install/nip05-dm-listener.service
[Unit]
Description=NIP-05 DM Listener — jorgenclaw.ai registration
After=network.target nostr-signer.service

[Service]
Type=simple
WorkingDirectory=%h/nostr-mcp-server/nip05-dm-listener
EnvironmentFile=%h/nostr-mcp-server/nip05-dm-listener/.env
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
