# nip05-dm-listener/install/nip05-dm-listener.service
[Unit]
Description=NIP-05 DM Listener — jorgenclaw.ai registration
After=network.target nostr-signer.service

[Service]
Type=simple
WorkingDirectory=%h/nostr-mcp-server/nip05-dm-listener
ExecStart=/home/jorgenclaw/.nvm/versions/node/v22.22.0/bin/node --env-file=.env index.js
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
